# crontab
0 0 * * * python /backup.py >> /var/log/cron.log 2>&1
